# templates/supervisor_netbox.conf.j2
[program:netbox]
command = /opt/netbox/bin/gunicorn -c /opt/netbox/gunicorn_config.py netbox.wsgi
directory = /opt/netbox/netbox/
user = www-data
